# macos/tart/scripts/com.gitea.actions.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gitea.actions.runner</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/admin/gitea-runner/gitea-runner.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/admin/Library/Logs/gitea-runner.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/admin/Library/Logs/gitea-runner.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/admin/.dotnet:/Users/admin/.dotnet/tools:/usr/local/share/dotnet:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>DOTNET_ROOT</key>
        <string>/Users/admin/.dotnet</string>
        <key>ANDROID_HOME</key>
        <string>/Users/admin/Library/Android/sdk</string>
        <key>ANDROID_SDK_ROOT</key>
        <string>/Users/admin/Library/Android/sdk</string>
    </dict>
</dict>
</plist>
